# net.gapul.keystats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>net.gapul.keystats</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Keystats.app 内から起動 = 入力監視の一覧にアプリのアイコン/名前で出る -->
    <string>__HOME__/Applications/Keystats.app/Contents/MacOS/keystatsd</string>
    <string>run</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <!-- launchd はシェルの env を継承しないので XDG をここで明示(CLIと参照先を揃える) -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>XDG_DATA_HOME</key><string>__DATA__</string>
    <key>XDG_STATE_HOME</key><string>__STATE__</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__STATE__/keystats/keystats.log</string>
  <key>StandardErrorPath</key>
  <string>__STATE__/keystats/keystats.log</string>
</dict>
</plist>
